# systemd/ai-telephone.service
[Unit]
Description=AI Telephone Backend
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 -m src.api_server
WorkingDirectory=/opt/ai-telephone
Restart=always
User=pi

[Install]
WantedBy=multi-user.target
